# etc/udev/rules.d/98-monitor.rules
KERNEL=="card0", SUBSYSTEM=="drm", ENV{DISPLAY}=":0", ENV{XAUTHORITY}="/home/markus/.Xauthority", RUN+="/usr/share/markus/hotplug-monitor.sh"
